# napic-files/files/etc/uci-defaults/93-console-password
#!/bin/sh
uci set system.@system[0].ttylogin='1'
uci commit system
