# napic-files/files/etc/uci-defaults/95-network
#!/bin/sh
uci set network.lan.device='eth0'
uci delete network.lan.type 2>/dev/null
uci commit network
