# napic-files/files/etc/uci-defaults/97-luci-theme
#!/bin/sh
uci set luci.main.mediaurlbase='/luci-static/openwrt-2020'
uci commit luci
