# napic-files/files/etc/uci-defaults/99-dhcp
#!/bin/sh
if [ ! -e /sys/class/net/eth1 ]; then
    # Одна сетевуха - не раздаём DHCP
    uci set dhcp.lan.ignore='1'
fi
uci commit dhcp
