# files/etc/uci-defaults/96-hostname
#!/bin/sh
uci set system.@system[0].hostname='napiwrt'
uci commit system
