# napi2-files/files/etc/uci-defaults/96-hostname
#!/bin/sh
uci set system.@system[0].hostname='napi2wrt'
uci commit system
