# napi2-files/files/etc/uci-defaults/97-luci-theme
#!/bin/sh
uci set luci.main.mediaurlbase='/luci-static/openwrt2020'
uci commit luci
